# project1/task_scheduler/config/crontab
5 * * * * root python3 monitor_node_health.py
