# project1/task_scheduler/config/crontab
5 * * * * root python3 /srv/src/monitor_node_health.py > /var/log/cron.log
