# project1/web_server/config/crontab
# This cron runs every minute
* * * * * python /srv/src/cron_job.py
